# overlay-root/usr/lib/systemd/user/laptop-backup.service
[Unit]
Description=Backup Documents and Pictures to local restic repo, then sync to Garage S3

[Service]
Type=oneshot
ExecStart=/usr/local/bin/laptop-backup
